# build/linux/slimrmm-agent.service
[Unit]
Description=SlimRMM Remote Monitoring & Management Agent
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/slimrmm-agent
WorkingDirectory=/var/lib/slimrmm
Restart=always
RestartSec=10
User=root
Group=root
Environment="SLIMRMM_SERVICE=1"

# Logging to system-typical location
StandardOutput=append:/var/log/slimrmm/agent.log
StandardError=append:/var/log/slimrmm/agent.log

[Install]
WantedBy=multi-user.target
