# build/macos/io.slimrmm.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.slimrmm.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/SlimRMM.app/Contents/MacOS/slimrmm-agent</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Library/Application Support/SlimRMM</string>
    <key>StandardOutPath</key>
    <string>/Library/Logs/SlimRMM/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/SlimRMM/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
