# build/macos/io.slimrmm.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.slimrmm.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/slimrmm-agent</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/var/lib/slimrmm</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/lib/slimrmm/log/agent.log</string>
    <key>StandardErrorPath</key>
    <string>/var/lib/slimrmm/log/agent.log</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
